# deploy/launchd/com.craig.skill.youtube.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.craig.skill.youtube</string>

    <!-- YouTube 요약봇 — 상시 리스너(monitor.py --listen):
         · 텔레그램 채널 명령(추가/삭제/list/run) 즉시 응답
         · 같은 프로세스가 schedule_interval_hours(config, 기본 6h)마다 새 영상 자동 감지·요약·전송·볼트로그
         getUpdates 를 이 프로세스 하나만 소유하므로 별도 주기잡과의 충돌이 없다. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/craigpark/anaconda3/bin/python3</string>
        <string>/Users/craigpark/Github/Craig-Skill/youtube-telegram-summary/monitor.py</string>
        <string>--listen</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/craigpark/Github/Craig-Skill/youtube-telegram-summary</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/Users/craigpark/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- 상시 구동: 로그인 시 시작 + 죽으면 자동 재시작 -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/craigpark/Github/Craig-Skill/logs/youtube.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/craigpark/Github/Craig-Skill/logs/youtube.err.log</string>
</dict>
</plist>
