# deploy/launchd/com.craig.skill.learn-retro.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.craig.skill.learn-retro</string>
  <key>ProgramArguments</key><array>
    <string>/Users/craigpark/anaconda3/bin/python3</string>
    <string>/Users/craigpark/Github/Craig-Skill/craig-telegram-study/pipeline/learn_retro.py</string>
    <string>--run</string></array>
  <key>WorkingDirectory</key><string>/Users/craigpark/Github/Craig-Skill/craig-telegram-study/pipeline</string>
  <key>EnvironmentVariables</key><dict><key>PATH</key><string>/opt/homebrew/bin:/Users/craigpark/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string><key>PYTHONUNBUFFERED</key><string>1</string></dict>
  <key>StartCalendarInterval</key><dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>/Users/craigpark/Github/Craig-Skill/logs/learn-retro.out.log</string>
  <key>StandardErrorPath</key><string>/Users/craigpark/Github/Craig-Skill/logs/learn-retro.err.log</string>
</dict></plist>
