# deploy/launchd/com.craig.skill.learn-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.craig.skill.learn-ingest</string>

    <!-- learn-ingest — incoming 큐 처리(catch-up). relay_bot 이 메시지마다 즉시 트리거하고,
         이 잡은 5분마다 밀린 큐를 처리하는 안전망. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/craigpark/anaconda3/bin/python3</string>
        <string>/Users/craigpark/Github/Craig-Skill/craig-telegram-study/pipeline/learn_ingest.py</string>
        <string>--queue</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/craigpark/Github/Craig-Skill/craig-telegram-study/pipeline</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/Users/craigpark/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/craigpark/Github/Craig-Skill/logs/learn-ingest.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/craigpark/Github/Craig-Skill/logs/learn-ingest.err.log</string>
</dict>
</plist>
